# replic-nagios.service
[Unit]
Description=Replication of the Nagios

[Service]
ExecStart=/etc/replic-nagios

[Install]
WantedBy=multi-user.target
